# service/oneminute/finish
#!/bin/sh

tasktxt=`cat /Users/stephen/.motd`
echo "stopped $tasktxt"
